# example_init_script.conf
description	"djangowake"

start on runlevel [2345]
stop on runlevel [!2345]

script
	cd "/srv/djangowake"
	exec sudo -u gunicorn gunicorn djangowake.wsgi --bind=127.0.0.1:8001
end script
